# scripts/services/disable-thp.service
[Unit]
Description=Disable Transparent Huge Pages (THP)
After=sysinit.target
Before=postgresql-bench.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo never > /sys/kernel/mm/transparent_hugepage/enabled && echo never > /sys/kernel/mm/transparent_hugepage/defrag'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
